# backend/systemd/dicom-watcher.service
[Unit]
Description=Dental DICOM Ingestion Watcher
After=network.target postgresql.service

[Service]
Type=simple
User=dicom
Group=dicom
WorkingDirectory=/opt/dental-portal/backend
EnvironmentFile=/opt/dental-portal/backend/.env
ExecStart=/opt/dental-portal/backend/.venv/bin/python -m app.services
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
